# ifconfig.service
[Unit]
Description = ifconfig web service

[Service]
ExecStart = /path/to/ifconfig/binary -mmdb /etc/GeoLite2-City.mmdb